# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Streamlit ────────────────────────────────────────────────
[program:streamlit]
command=streamlit run app.py
    --server.port=%(ENV_PORT)s
    --server.address=0.0.0.0
    --server.headless=true
    --server.fileWatcherType=none
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=HOME="/root",PATH="/usr/local/bin:%(ENV_PATH)s"

; ── Background worker ────────────────────────────────────────
[program:worker]
command=python worker.py
directory=/app
autostart=true
autorestart=true
startretries=5
; Give Streamlit 5 s to init DB before worker starts
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=HOME="/root",PATH="/usr/local/bin:%(ENV_PATH)s"

; ── Process group ────────────────────────────────────────────
[group:app]
programs=streamlit,worker
